# docker/evelyn-native/systemd/parquet-writer.service
[Unit]
Description=Parquet Writer (Evelyn)
After=network.target flashmq.service
Wants=flashmq.service

[Service]
Type=simple
ExecStart=/usr/local/bin/parquet-writer --config /etc/evelyn/writer.yaml
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
