# scripts/writer-compare-test.service
# writer-compare-test.service
# Schema comparison benchmark: norm-long vs long-cmp vs wide-pivot
#
# Install:
#   sudo cp scripts/writer-compare-test.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#
# Run:
#   sudo systemctl start writer-compare-test
#   journalctl -u writer-compare-test -f
#
# On .34 the output lands on /mnt/tort-sdi/bench/
# Override OUTDIR/SWEEPS via drop-in:
#   sudo systemctl edit writer-compare-test
#   [Service]
#   Environment=SWEEPS=300
#   Environment=OUTDIR=/mnt/tort-sdf/bench

[Unit]
Description=Parquet Writer Schema Comparison Benchmark
After=network.target flashmq.service
Wants=flashmq.service

[Service]
Type=oneshot
RemainAfterExit=no
User=phil

# Auto-detect repo root from service file location
# Works for both ~/data_server (.34) and ~/work/gen-ai/data_server (.46)
ExecStart=/bin/bash -c '\
  REPO=$(find /home/phil -maxdepth 4 -name "writer_compare_test.sh" \
         -path "*/scripts/*" 2>/dev/null | head -1 | xargs dirname | xargs dirname); \
  echo "repo=$REPO"; \
  exec bash "$REPO/scripts/writer_compare_test.sh" \
    --sweeps ${SWEEPS:-600} \
    --outdir ${OUTDIR:-} \
    --site   ${SITE:-SITE_A} \
    ${COV:-}'

Environment=SWEEPS=600
Environment=SITE=SITE_A
Environment=COV=
# Environment=OUTDIR=/mnt/tort-sdi/bench   # uncomment to override

# Give the test plenty of time (600 sweeps × 1s + analysis ≈ 15min)
TimeoutStartSec=1200

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
